# flags.service
[Unit]
Description=A Flags server daemon

[Service]
ExecStart=/usr/bin/flags
User=flags

[Install]
WantedBy=network.target
